# newsboat/.config/systemd/user/newsboat.service
[Unit]
Description=Refreshes and notifies for newsboat

[Service]
ExecStart=newsboat -x reload

[Install]
WantedBy=default.target
