# newsboat/.config/systemd/user/newsboat.timer
[Unit]
Description=Update Newsboat every hour

[Timer]
OnBootSec=1hour
OnActiveSec=1hour
OnUnitActiveSec=1hour

[Install]
WantedBy=timers.target
