# davmail/.config/systemd/user/checkemails@inbox.timer
[Unit]
Description=Check Inbox Every Minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
