# systemd/HOME/.config/systemd/user/background.timer
[Unit]
Description=Change Wallpaper every minute

[Timer]
OnBootSec=1min
OnActiveSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
