# systemd/HOME/.config/systemd/user/checkemails@.service
[Unit]
Description=Check emails for %i

[Service]
ExecStart=/home/jonathan/.bin/emails/check-emails %i

[Install]
WantedBy=default.target
