# systemd/HOME/.config/systemd/user/checkemails@inbox.timer
[Unit]
Description=Check Inbox Every Minute

[Timer]
OnCalendar=Mon..Fri 9..18:*:00

[Install]
WantedBy=timers.target
