# systemd/HOME/.config/systemd/user/lookAway.service
[Unit]
Description=Look Away

[Service]
ExecStart=/home/jonathan/.config/systemd/user/lookAway

[Install]
WantedBy=default.target
