# systemd/.config/systemd/user/background.service
[Unit]
Description=Background Switcher

[Service]
ExecStart=/home/jonathan/.config/systemd/user/randomBackground

[Install]
WantedBy=default.target
