# systemd/.config/systemd/user/stock.service
[Unit]
Description=Check stock

[Service]
ExecStart=/home/jonathan/.config/systemd/user/checkStock

[Install]
WantedBy=default.target
